# systemd/.config/systemd/user/tldr.service
[Unit]
Description=Update the TLDR cache

[Service]
ExecStart=/usr/bin/tldr --update
Type=oneshot
