# touchegg/.config/systemd/user/touchegg-client.service
[Unit]
Description=Touchegg client
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/touchegg

[Install]
WantedBy=graphical.target
